# systemd/plexmediaserver.service.d/override.conf
[Unit]
After=media.service
Requires=media.service rclone-mma.service rclone-movies.service rclone-tv.service

[Service]
User=felix
Group=felix
